# build/conf/supervisord.conf
[supervisord]
logfile=/xhgui/log/supervisor.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --bind_ip_all --dbpath /xhgui/mongodb --port 27017
stdout_logfile=/xhgui/log/mongodb.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/xhgui/log/mongodb.err.log
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_logfile=/xhgui/log/nginx.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/xhgui/log/nginx.err.log
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/local/sbin/php-fpm
stdout_logfile=/xhgui/log/php-fpm.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/xhgui/log/php-fpm.err.log
stderr_logfile_maxbytes=0

[program:start]
autorestart=false
autostart=true
command=/frx/start
exitcodes=0
startretries=0
stdout_logfile=/xhgui/log/start.out.log
stdout_logfile_maxbytes=0
stderr_logfile=/xhgui/log/start.err.log
stderr_logfile_maxbytes=0
